# csdcinstaller/rc.local.service

[Install]  
WantedBy=multi-user.target  
Alias=rc-local.service

